# colorblind-probe.service
[Unit]
Description=Colorblind Probe HTTP Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/colorblind-probe
ExecStart=/usr/bin/python3 -m http.server 8000 --bind 0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
